# bitsy.service
[Unit]
Description=Bitsy - AI Robot Car
After=network.target sound.target
Wants=network.target

[Service]
Type=forking
User=jmunning
Group=jmunning
WorkingDirectory=/home/jmunning
ExecStart=/home/jmunning/start_bitsy.sh
PIDFile=/home/jmunning/bitsy.pid
Restart=always
RestartSec=10
StandardOutput=append:/home/jmunning/bitsy_service.log
StandardError=append:/home/jmunning/bitsy_service.log

# Hardware access permissions
SupplementaryGroups=gpio i2c spi audio

# Environment variables
Environment=PATH=/home/jmunning/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=JACK_NO_START_SERVER=1
Environment=JACK_NO_AUDIO_RESERVATION=1

[Install]
WantedBy=multi-user.target 
