# bitsy.service
[Unit]
Description=Bitsy - AI Robot Car
After=network.target sound.target pulseaudio.service
Requires=pulseaudio.service
Wants=network.target

[Service]
Type=forking
User=jmunning
Group=jmunning
WorkingDirectory=/home/jmunning
Environment="XDG_RUNTIME_DIR=/run/user/1000"
Environment="PULSE_RUNTIME_PATH=/run/user/1000/pulse"
Environment="PULSE_SERVER=unix:/run/user/1000/pulse/native"
ExecStartPre=/bin/bash -c 'pulseaudio --start --log-target=syslog || true'
ExecStart=/bin/bash -c 'source /home/pi/venv/bin/activate && python3 /home/pi/chatgpt_with_driving.py'
PIDFile=/home/jmunning/bitsy.pid
Restart=always
RestartSec=10
StandardOutput=append:/home/jmunning/bitsy_service.log
StandardError=append:/home/jmunning/bitsy_service.log

# Hardware access permissions
SupplementaryGroups=gpio i2c spi audio

# Environment variables
Environment=PATH=/home/jmunning/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=JACK_NO_START_SERVER=1
Environment=JACK_NO_AUDIO_RESERVATION=1

[Install]
WantedBy=multi-user.target 
